# systemd/recordo.service
[Unit]
Description=Recordo daemon (gravador de reuniões fricção-zero / TDAH-friendly)
Documentation=https://github.com/caio-hat/recordo
After=pipewire.service pipewire-pulse.service graphical-session.target
Wants=pipewire.service

[Service]
Type=simple
# __RECORDO_BIN__ é substituído pelo setup.sh (path absoluto do entrypoint no venv)
ExecStart=__RECORDO_BIN__ --daemon
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:%h/.local/bin:/home/linuxbrew/.linuxbrew/bin
# v0.2.3: forçar locale UTF-8 (systemd user units herdam locale C/POSIX por
# default, causando UnicodeDecodeError em subprocess.run(text=True) quando
# stdout/stderr tem chars multibyte)
Environment=LC_ALL=C.UTF-8
Environment=LANG=C.UTF-8
Environment=PYTHONIOENCODING=utf-8
StandardInput=null
StandardOutput=append:%t/recordo.daemon.log
StandardError=append:%t/recordo.daemon.log

[Install]
WantedBy=default.target
